# plugins/oracle/deploy/oci-migrator.service
# /etc/systemd/system/oci-migrator.service
# Systemd unit for the OCI preemptible live-migration daemon.
# Install on every build instance in your CI/CD fleet.

[Unit]
Description=OCI Preemptible Live Migrator
Documentation=https://github.com/bons/bons-ci/plugins/oracle
# Start after the shared migration volume is mounted.
After=network-online.target mnt-migration.mount
Wants=network-online.target mnt-migration.mount
# Restart aggressively — if this daemon dies, migration cannot trigger.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/oci-migrator
Restart=always
RestartSec=2s

# Run as root — CRIU requires CAP_SYS_PTRACE and CAP_SYS_ADMIN.
User=root
Group=root

# Capabilities explicitly required:
#   CAP_SYS_PTRACE  — CRIU ptrace-based process introspection
#   CAP_SYS_ADMIN   — namespace operations, cgroup freezer
#   CAP_NET_ADMIN   — netlink route/address manipulation
#   CAP_DAC_READ_SEARCH — reading /proc/<pid>/maps etc.
AmbientCapabilities=CAP_SYS_PTRACE CAP_SYS_ADMIN CAP_NET_ADMIN CAP_DAC_READ_SEARCH

# Security hardening (we restrict what we don't need).
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

# Environment.
EnvironmentFile=-/etc/oci-migrator/env
Environment=OCI_MIGRATOR_OCI_METADATA_ENDPOINT=http://169.254.169.254/opc/v2

# Logging via journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oci-migrator

# Give the daemon enough time to complete an in-flight migration.
TimeoutStopSec=130s
KillSignal=SIGTERM
KillMode=process

[Install]
WantedBy=multi-user.target
